# crontab
0 8 * * * cd /home/ec2-user/houseingLottery && python3 house.py >> log.txt 2>&1
